# target/etc/systemd/system/systemd-networkd-wait-online.service.d/override.conf
[Unit]
ConditionPathExistsGlob=/run/systemd/network/*
