# build/package/target/usr/lib/systemd/system/fwait-bootfile.service
[Unit]
Description=Wait for bootfile
Before=inspector.service

[Service]
Type=notify
NotifyAccess=all
Restart=no
EnvironmentFile=/etc/default/inspector
ExecStart=/usr/bin/fwait --file ${BOOTFILE_NAME}
