# build/package/target/usr/lib/systemd/system/runner-pre.target

[Unit]
Description=Preparation for Runner
Requires=basic.target
After=basic.target
AllowIsolate=yes
