# build/package/target/usr/lib/systemd/system/shell-pre.target

[Unit]
Description=Preparation for Shell
Requires=basic.target
After=basic.target
AllowIsolate=yes
