# build/package/target/usr/lib/systemd/system/shell.service
[Unit]
Description=Shell
Requires=shell-pre.target
After=shell-pre.target
IgnoreOnIsolate=yes
ConditionPathExists=/dev/tty0

[Service]
Type=idle
Environment=TERM=linux HOME=/root PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/local/bin/shell
Restart=no
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
TTYReset=no
TTYVHangup=no
TTYVTDisallocate=no
