# target/usr/lib/systemd/system/bootstrap@.service
[Unit]
Description=Bootstrap on %I
IgnoreOnIsolate=yes
OnFailure=KLoader-rescue.target
OnSuccessJobMode=replace-irreversibly

[Service]
Type=oneshot
Environment=TERM=linux HOME=/
ExecStart=/usr/local/bin/bootstrap
Restart=no
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/%I
TTYReset=no
TTYVHangup=no
TTYVTDisallocate=no
IgnoreSIGPIPE=no
SendSIGHUP=no

[Install]
WantedBy=KLoader.target
DefaultInstance=console
