# target/usr/lib/systemd/system/bootstrap.service
[Unit]
Description=Firestarter Bootstrap
After=sysinit.target
Before=systemd-networkd.service

[Service]
Type=oneshot
Restart=no
ExecStart=/opt/firestarter/bootstrap
