# target/usr/lib/systemd/system/bootstrap.service
[Unit]
Description=Firestarter Bootstrap
After=sysinit.target
Before=systemd-networkd.service

[Service]
Type=oneshot
Restart=no
ExecStartPre=-/usr/bin/mkdir /run/firestarter /run/firestarter/bootstrap /run/systemd/network
ExecStart=/opt/firestarter/bootstrap
