# target/usr/lib/systemd/system/bootstrap.service
[Unit]
Description=Bootstrap
After=sysinit.target
Before=systemd-networkd.service

[Service]
Type=oneshot
Restart=no
ExecStartPre=/usr/bin/mkdir -p /run/bootstrap /run/systemd/network
ExecStart=/usr/local/bin/bootstrap
