# target/usr/lib/systemd/system/firestarter-pre.target

[Unit]
Description=Preparation for Firestarter
Requires=basic.target
After=basic.target
AllowIsolate=yes
