# target/usr/lib/systemd/system/firestarter.target

[Unit]
Description=Firestarter
Requires=basic.target
After=basic.target
AllowIsolate=yes
