# target/usr/lib/systemd/system/KLoader.target
[Unit]
Description=KLoader
AssertPathExists=/etc/initrd-release
Requires=basic.target
After=basic.target
AllowIsolate=yes
