# target/usr/lib/systemd/system/KLoader-rescue.service
[Unit]
Description=KLoader Rescue

[Service]
Type=oneshot
Environment=TERM=linux HOME=/
ExecStart=/opt/KLoader/rescue
Restart=no
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/console
TTYReset=no
TTYVHangup=no
TTYVTDisallocate=no
