# target/usr/lib/systemd/system/mounter.service
[Unit]
Description=Mounter
After=sysinit.target bootstrap.service
ConditionPathExists=/run/bootstrap/PartitionUUID.env

[Service]
Type=oneshot
Restart=no
ExecStartPre=/usr/bin/mkdir -p /efi
ExecStart=/usr/local/bin/mounter
EnvironmentFile=/run/bootstrap/PartitionUUID.env
